# telegram-clothing-order-bot.service
[Unit]
Description=Telegram Clothing Order Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=phyoshan8
WorkingDirectory=/home/phyoshan8/automation-projects/telegram-clothing-order-bot
EnvironmentFile=/home/phyoshan8/automation-projects/telegram-clothing-order-bot/.env
ExecStart=/home/phyoshan8/automation-projects/telegram-clothing-order-bot/.venv/bin/python -m src.bot
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
